# k8s-drain-deb/DEBIAN/postinst
#!/bin/bash

# Перезагрузите systemd и активируйте сервис
systemctl daemon-reload
systemctl enable k8s-drain.service
